# subway_train_times/subway-display.service
[Unit]
Description=Subway Display Image Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-zero/subway_train_times
ExecStart=/usr/bin/python3 /home/pi/pi-zero/subway_train_times/subway_server.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target